# extIP.service
[Unit]
Description=extIP

[Service]
Type=simple
ExecStart=/extIP.sh
